# etc/udev/rules.d/70-automount.rules
ACTION=="add", KERNEL=="sd?1", SUBSYSTEM=="block", RUN+="/bin/systemctl --no-block start backup@%k.service"
