# etc/udev/rules.d/70-automount.rules
ACTION=="add", KERNEL=="sd*", SUBSYSTEM=="block", ENV{ID_FS_TYPE}=="vfat", RUN+="/bin/systemctl --no-block start backup@%k.service"
