# etc/systemd/system/backup@.service
[Unit]
Description = Backup to USB Flash Drive
BindsTo = dev-%i.device

[Service]
Type = simple
ExecStart = /home/pi/run-copy.sh %I
